# backend/deploy/hf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
autorestart=true
startsecs=10
startretries=3
stopsignal=INT
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:detection]
command=uvicorn services.detection.main:app --host 0.0.0.0 --port 8001 --log-level warning
directory=/app/backend
environment=PYTHONPATH="/app/backend:/app/backend/services/detection",SERVICE_NAME="detection",DATABASE_URL="postgresql+psycopg2://parksight:parksight@localhost:5432/parksight",REDIS_URL="",KAFKA_BOOTSTRAP_SERVERS="",JWT_SECRET_KEY="hf-demo-jwt-secret-key",CORS_ORIGINS="*",RATE_LIMIT_ENABLED="false"
user=root
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
[program:scoring]
command=uvicorn services.scoring.main:app --host 0.0.0.0 --port 8002 --log-level warning
directory=/app/backend
environment=PYTHONPATH="/app/backend",SERVICE_NAME="scoring",DATABASE_URL="postgresql+psycopg2://parksight:parksight@localhost:5432/parksight",REDIS_URL="",KAFKA_BOOTSTRAP_SERVERS="",JWT_SECRET_KEY="hf-demo-jwt-secret-key",CORS_ORIGINS="*",RATE_LIMIT_ENABLED="false"
user=root
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# no depends_on — supervisor doesn't support it

[program:dispatch]
command=uvicorn services.dispatch.main:app --host 0.0.0.0 --port 8003 --log-level warning
directory=/app/backend
environment=PYTHONPATH="/app/backend",SERVICE_NAME="dispatch",DATABASE_URL="postgresql+psycopg2://parksight:parksight@localhost:5432/parksight",REDIS_URL="",KAFKA_BOOTSTRAP_SERVERS="",JWT_SECRET_KEY="hf-demo-jwt-secret-key",CORS_ORIGINS="*",RATE_LIMIT_ENABLED="false"
user=root
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# no depends_on — supervisor doesn't support it

[program:analytics]
command=uvicorn services.analytics.main:app --host 0.0.0.0 --port 8004 --log-level warning
directory=/app/backend
environment=PYTHONPATH="/app/backend",SERVICE_NAME="analytics",DATABASE_URL="postgresql+psycopg2://parksight:parksight@localhost:5432/parksight",REDIS_URL="",KAFKA_BOOTSTRAP_SERVERS="",JWT_SECRET_KEY="hf-demo-jwt-secret-key",CORS_ORIGINS="*",RATE_LIMIT_ENABLED="false"
user=root
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# no depends_on — supervisor doesn't support it

[program:officers]
command=uvicorn services.officers.main:app --host 0.0.0.0 --port 8005 --log-level warning
directory=/app/backend
environment=PYTHONPATH="/app/backend",SERVICE_NAME="officers",DATABASE_URL="postgresql+psycopg2://parksight:parksight@localhost:5432/parksight",REDIS_URL="",KAFKA_BOOTSTRAP_SERVERS="",JWT_SECRET_KEY="hf-demo-jwt-secret-key",CORS_ORIGINS="*",RATE_LIMIT_ENABLED="false"
user=root
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# no depends_on — supervisor doesn't support it

[program:alerts]
command=uvicorn services.alerts.main:app --host 0.0.0.0 --port 8006 --log-level warning
directory=/app/backend
environment=PYTHONPATH="/app/backend",SERVICE_NAME="alerts",DATABASE_URL="postgresql+psycopg2://parksight:parksight@localhost:5432/parksight",REDIS_URL="",KAFKA_BOOTSTRAP_SERVERS="",JWT_SECRET_KEY="hf-demo-jwt-secret-key",CORS_ORIGINS="*",RATE_LIMIT_ENABLED="false"
user=root
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# no depends_on — supervisor doesn't support it

[program:nginx]
command=nginx -g "daemon off;"
user=root
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# no depends_on — supervisor doesn't support it
